# safaodev.service
[Unit]
Description=My Service

[Service]
Type=simple
ExecStart= /home/safa/bilisim_script.sh

[Install]
WantedBy=multi-user.target
